# deploy/macos/com.hermes.mobile-remote.dashboard.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.hermes.mobile-remote.dashboard</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>{{BIN_DIR}}/run_dashboard.sh</string>
  </array>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>StandardOutPath</key>
  <string>{{RUNTIME_DIR}}/dashboard.out.log</string>
  <key>StandardErrorPath</key>
  <string>{{RUNTIME_DIR}}/dashboard.err.log</string>
</dict>
</plist>
